# deploy/com.landline.daemon-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Landline watchdog template. Runs every 5 minutes; if the daemon has
  disappeared from launchd entirely (bootout without bootstrap, boot glitch,
  crash-loop giving up), re-bootstraps it. This is orthogonal to
  KeepAlive on the daemon plist — KeepAlive restarts a *crashed* process
  in place; the watchdog rescues the *label* if it fell off launchd.

  Copy to ~/Library/LaunchAgents/, then edit:
    - ProgramArguments => absolute path to deploy/watchdog.sh in your checkout
    - StandardOutPath / StandardErrorPath => wherever the watchdog log goes

  The Label and the LANDLINE_PLIST / LANDLINE_LABEL env values here match
  deploy/watchdog.sh's defaults and deploy/com.landline.telegram-daemon.plist's
  label, so a user following docs/SETUP.md verbatim can load this with no
  extra edits beyond the /Users/YOU paths.

  Environment overrides the watchdog script honors (see deploy/watchdog.sh):
    LANDLINE_PLIST              path to the main daemon plist
    LANDLINE_LABEL              launchd label to re-bootstrap
    LANDLINE_WATCHDOG_LOG       log file the watchdog appends to
    LANDLINE_KEYCHAIN_ACCOUNT   Keychain account (default: landline)

  Load with:
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.landline.daemon-watchdog.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.landline.daemon-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOU/claude-landline/deploy/watchdog.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LANDLINE_PLIST</key>
        <string>/Users/YOU/Library/LaunchAgents/com.landline.telegram-daemon.plist</string>
        <key>LANDLINE_LABEL</key>
        <string>com.landline.telegram-daemon</string>
        <key>LANDLINE_KEYCHAIN_ACCOUNT</key>
        <string>landline</string>
    </dict>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/YOU/.landline/logs/telegram-daemon/watchdog.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/.landline/logs/telegram-daemon/watchdog.err.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
